# docker/files/etc/services.d/strobe/run
#!/usr/bin/with-contenv sh

exec 2>&1

cd /strobe/

exec ./apps/elvis/bin/run-production
